# deploy/main/orasage-main.service
[Unit]
Description=OraSage Main Portal (Next.js)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/orasage/main
Environment=NODE_ENV=production
Environment=PORT=3100
Environment=HOSTNAME=127.0.0.1
ExecStart=/usr/local/bin/npm start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
