# deploy/admin/orasage-admin.service
[Unit]
Description=OraSage Admin Portal (Next.js)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/orasage/admin
EnvironmentFile=-/opt/orasage/.env
Environment=NODE_ENV=production
Environment=PORT=3103
Environment=HOSTNAME=127.0.0.1
ExecStart=/usr/local/bin/npm start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
