# deploy/cms/orasage-cms.service
[Unit]
Description=OraSage CMS (Payload)
After=network.target postgresql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/orasage/cms
EnvironmentFile=-/opt/orasage/cms/.env
Environment=NODE_ENV=production
Environment=PORT=3120
Environment=HOSTNAME=127.0.0.1
Environment=CMS_MEDIA_DIR=/var/lib/orasage/cms-media
ExecStart=/usr/local/bin/npm start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
